# ci/service/gitlab-runner
#!/sbin/openrc-run
description="Gitlab-runner service"

depend() {
    need net
}

start()  {
    ebegin "Starting gitlab-runner"
    start-stop-daemon --start --exec /usr/bin/gitlab-runner -- "run" "--user=gitlab-runner" "--working-directory=/home/gitlab-runner" "--config=/etc/gitlab-runner/config.toml" "--syslog"
    eend $?
    #gitlab-runner run --user=gitlab-runner --working-directory=/home/gitlab-runner --syslog 
}
